# gluster-pv.service
[Unit]
After=glusterd.service

[Service]
ExecStart=/srv/scripts/gluster-pv.sh

